# 0-mysqlbackup.cron
#!/bin/sh

# Password for MySQL root
password=

target_dir="/var/lib/mysql/bak"

mysqldump \
  --user=root \
  --password=$password \
  --lock-all-tables \
  --flush-logs \
  --default-character-set=utf8 \
  --all-databases \
  > $target_dir/mysql.dump.new

if [ $? != 0 ]; then
  echo "[ERROR] mysqldump was failed." 1>&2
  exit 1
fi

mv -f $target_dir/mysql.dump.new $target_dir/mysql.dump
