# 0-mysqlbackup.cron
#!/bin/sh

# Password for MySQL root
password=

mysql_data_dir='/var/lib/mysql'
binlog_dir="$mysql_data_dir/binlog"
backup_dir="$mysql_data_dir/bak"

mysqldump \
--user=root \
--password=$password \
--lock-all-tables \
--flush-logs \
--default-character-set=utf8 \
--all-databases \
> $backup_dir/mysql.dump.new

if [ $? != 0 ]; then
  echo "[ERROR] mysqldump was failed." 1>&2
  exit 1
fi

mv -f $backup_dir/mysql.dump.new $backup_dir/mysql.dump
